# plugins/personas-mesh/systemd/personas-mesh-windows-bridge.service
[Unit]
Description=personas-mesh Windows-FS bridge — sync /mnt/c/Users/*/.personas/*
Documentation=https://github.com/kickinrad/personas
# DISABLED BY DESIGN — do NOT enable. `personas-mesh-wsl` is the single
# canonical Layer-1 git-mesh path. Pre-Phase-4 (~/.personas was a symlink to
# /mnt/c/...) this bridge would have churned the same inode; post-Phase-4 the
# two are separate dirs holding the SAME personas, so running this alongside
# wsl dual-pushes identical repos to the hub and causes merge conflicts.
# Windows-side changes reach the mesh via personas-mesh-user-sync (Layer-2
# rsync) -> ~/.personas -> wsl. Re-enable ONLY if the Windows dir ever becomes
# a genuinely independent node (different personas than ~/.personas).

[Service]
Type=oneshot
Nice=10
ExecStart=%h/.local/bin/personas-mesh-sync-all /mnt/c/Users/wilst/.personas
StandardOutput=journal
StandardError=journal
