# plugins/personas-mesh/systemd/personas-mesh-windows-bridge.timer
[Unit]
Description=personas-mesh Windows-FS bridge timer — runs every 3 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=3min
AccuracySec=15s
Persistent=true

[Install]
WantedBy=timers.target
