# plugins/personas-mesh/systemd/personas-mesh-github-mirror.timer
[Unit]
Description=personas-mesh GitHub mirror timer — hourly async backup push

[Timer]
OnBootSec=15min
OnUnitActiveSec=1h
AccuracySec=5min
Persistent=true

[Install]
WantedBy=timers.target
